# deploy/cpamc-sidecar.service
[Unit]
Description=CPAMC Sidecar Manager
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=cpamc
Group=cpamc
WorkingDirectory=/opt/cpamc-sidecar
Environment=NODE_ENV=production
EnvironmentFile=/etc/cpamc-sidecar/cpamc-sidecar.env
ExecStart=/usr/bin/node /opt/cpamc-sidecar/server.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/cpamc-sidecar

[Install]
WantedBy=multi-user.target
